# dist/supervisord.conf
[supervisord]
nodaemon=true

[program:smtp-sink]
command=/usr/sbin/smtp-sink -u nobody -d "/opt/spamsink/mails/dincoming-%%Y%%m%%d%%H%%M." -c -h %(ENV_SINK_HOSTNAME)s 0.0.0.0:%(ENV_SINK_PORT)s 10
autorestart=true
